# etc/r2d8.crontab
# m h  dom mon dow   command
# run the supervisor if it is not 
*/1 * * * * /home/r2d8/r2d8-scripts/r2_supervisor.sh >> /home/r2d8/r2d8-scripts/logs/r2d8.log 2>&1

